# units/alc-service@.service.tmpl
[Unit]
# alc-service@<binary-name>.service
# 例: alc-service@rust-alc-api.service / alc-service@gateway.service
# /target/release/<binary-name> を起動する。/target は cargo target volume が
# mount される場所 (build.sh が CARGO_TARGET_DIR=/target で書き込む)。
Description=alc dev service %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/alc-service@%i.env
WorkingDirectory=/src
# binary が無いと ExecStart で fail するので、user が build.sh で先に作る前提。
ExecStartPre=/usr/bin/test -x /target/release/%i
ExecStart=/target/release/%i
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
